# attendance
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="attendance"
description="Attendance zkteco"
command="python3 /root/Zkteco/app.py"